# systemd/airplay-agent.service
[Unit]
Description=AirPlay Agent MCP Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ollama
WorkingDirectory=/home/ollama/airplay-agent
Environment="PATH=/home/ollama/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/ollama/venv/lib/python3.11/site-packages"
ExecStart=/home/ollama/venv/bin/python -m airplay_agent.mcp_server
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
