# systemd/castmasta.service
[Unit]
Description=CastMasta MCP Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ollama
WorkingDirectory=/home/ollama/castmasta
Environment="PATH=/home/ollama/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/ollama/venv/bin/python -m castmasta.mcp_server
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/ollama/.castmasta
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictRealtime=yes
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target
