# systemd/castmasta@.service
[Unit]
Description=CastMasta MCP Server (per-user)
After=network.target

[Service]
Type=simple
Environment="PATH=/home/ollama/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/ollama/venv/bin/python -m castmasta.mcp_server
Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=%h/.castmasta
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictRealtime=yes
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service

[Install]
WantedBy=default.target
