# debian/castmasta.service
[Unit]
Description=CastMasta AirPlay/Google Cast MCP agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=castmasta
Group=castmasta
ExecStart=/usr/lib/castmasta/bin/castmasta-mcp --http --host 0.0.0.0 --port 16384
Restart=on-failure
RestartSec=5
StateDirectory=castmasta
StateDirectoryMode=0700
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/castmasta
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
